# php-publication/supervisor/publication-data-clean.ini
[program:short_url]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/xingstation-publication/artisan queue:work --queue=data-clean --tries=1
autostart=true
autorestart=true
numprocs=2
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/xingstation-publication/storage/logs/queue.log