# php/supervisor/jingsaas-wechat-outh.ini
[program:face]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work --queue=wechat-outh --tries=1
autostart=true
autorestart=true
numprocs=2
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/queue.log